# deployment/community_gunicorn_supervisor.conf
[program:demo] 
command = /var/www/html/pythonproject/django-community/deployment/community_gunicorn_service         ; Command to start app 
user = kawsar                                                                                     ; User to run as 
stdout_logfile = /var/www/html/pythonproject/django-community/logs/gunicorn_supervisor.log  ; Where to write log messages 
redirect_stderr = true                                                                        ; Save stderr in the same log 
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8 